# recipes-hello-world/hello-world/hello-world.service
#
# Hello-world systemd unit file
#
[Unit]
Description="Hello, world"

[Service]
ExecStart=/usr/bin/hello &
